# usr/lib/systemd/user/app-org.gnome.Software@autostart.service.d/moeOS-rc.conf
[Service]
CPUQuota=30%
CPUWeight=idle
IOWeight=1