# usr/lib/systemd/user/org.gnome.Shell@user.service.d/moeOS-override.service
[Service]
CPUWeight=5000
MemorySwapMax=0