# usr/lib/systemd/system/bluetooth.service.d/moeOS-override.conf
[Service]
Nice=-5
CPUWeight=150
IOWeight=200
IOSchedulingClass=realtime