# usr/lib/systemd/system/bluetooth.service.d/moeOS-override.conf
[Service]
IOWeight=200
IOSchedulingClass=realtime