# usr/lib/systemd/system/gdm.service.d/moeOS-override.conf
[Unit]
Wants=bluetooth.service rtkit-daemon.service
After=bluetooth.service rtkit-daemon.service
Requires=rtkit-daemon.service