# usr/lib/systemd/system/snapper-timeline.service.d/moeOS-rc.conf
[Unit]
ConditionACPower=true

[Service]
MemoryHigh=20%
CPUQuota=10%
CPUWeight=idle
IOWeight=1
IOReadBandwidthMax=/ 50M