# usr/lib/systemd/system/cups-browsed.service.d/moeOS-hardening.conf
[Service]
ProtectKernelTunables=yes
LockPersonality=yes
ProtectHostname=yes
ProtectKernelModules=yes
RestrictRealtime=yes
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid
PrivateUsers=yes
UMask=077
DynamicUser=yes
ProtectClock=yes
ProtectKernelLogs=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
IPAddressAllow=localhost link-local
CPUWeight=20