# usr/lib/systemd/system/chronyd.service.d/moeOS.conf
[Service]
BindPaths=/usr/share/moeOS-Docs/chrony.conf:/etc/chrony.conf
