# usr/lib/systemd/user/xdg-desktop-portal-gnome.service.d/1-moe.conf
[Service]
Environment=GDK_DEBUG=gl-no-fractional
Environment=GSK_RENDERER=gl
