# usr/lib/systemd/system/nvidia-suspend.service.d/moeOS-nvidia-sleep.conf
[Install]
# This broke systemctl sleep
RequiredBy=systemd-suspend-then-hibernate.service
RequiredBy=systemd-hibernate.service
RequiredBy=systemd-suspend.service
RequiredBy=systemd-hybrid-sleep.service

[Unit]
Before=systemd-hibernate.service
Before=systemd-hybrid-sleep.service
Before=systemd-suspend-then-hibernate.service
Before=systemd-suspend.service
