# usr/lib/systemd/system/nvidia-suspend.service.d/moeOS-nvidia-sleep.conf
[Install]
# This broke systemctl sleep
WantedBy=systemd-suspend-then-hibernate.service
