# usr/lib/systemd/system/nvidia-suspend.service.d/moeOS-workaround.conf
[Unit]
Before=sleep.target

[Install]
RequiredBy=sleep.target