# usr/lib/systemd/user/kde-baloo.service.d/moeOS-override.conf
[Service]
CPUWeight=idle
IOWeight=10
CPUQuota=30%
MemoryHigh=500M
