# usr/lib/systemd/user/kde-baloo.service.d/moeOS-override.conf
[Service]
CPUWeight=50
#IOWeight=10
CPUQuota=50%
MemoryHigh=500M
