# usr/lib/systemd/user/kde-baloo.service.d/moeOS-override.conf
[Service]
CPUWeight=idle
IOWeight=1
CPUQuota=20%
MemoryHigh=1G
