# usr/lib/systemd/system/nvidia-hibernate.service.d/moeOS-override.conf
[Service]
ExecStart=
ExecStart=bash /usr/lib/moeOS/nvsuspend