# usr/lib/systemd/system/nvidia-hibernate.service.d/moeOS-nvidia-sleep.conf
[Install]
RequiredBy=systemd-hibernate.service

[Unit]
Before=systemd-hibernate.service
