# usr/lib/systemd/system/NetworkManager-wait-online.service.d/moeOS.conf
[Service]
Environment=NM_ONLINE_TIMEOUT=40
Restart=on-failure
RestartSec=5s
RestartMode=direct