# usr/lib/systemd/system/NetworkManager-wait-online.service.d/moeOS.conf
[Unit]
Wants=moeOS-proxy.service
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Slice=moeOS-noproxy.slice
ExecStartPre=!resolvectl flush-caches
Environment=NM_ONLINE_TIMEOUT=50
Restart=on-failure
RestartSec=3s
RestartMode=direct