# usr/lib/systemd/system/NetworkManager-wait-online.service.d/moeOS.conf
[Unit]
Wants=moeOS-proxy.service

[Service]
Environment=NM_ONLINE_TIMEOUT=120
Restart=on-failure
RestartSec=5s
#RestartMode=direct