# usr/lib/systemd/system/NetworkManager-wait-online.service.d/moeOS.conf
[Unit]
Conflicts=sleep.target
After=systemd-suspend.service
After=systemd-hibernate.service
After=systemd-suspend-then-hibernate.service

[Service]
Environment=NM_ONLINE_TIMEOUT=120
Restart=on-failure
RestartSec=5s
#RestartMode=direct

[Install]
WantedBy=systemd-suspend.service
WantedBy=systemd-hibernate.service
WantedBy=systemd-suspend-then-hibernate.service