# usr/lib/systemd/system/user@.service.d/10-moe-oomd-per-slice-defaults.conf
[Service]
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=60%
MemoryPressureWatch=auto
