# usr/lib/systemd/system/user@.service.d/moeOS.conf
[Slice]
CPUWeight=10000

[Unit]
Wants=bluetooth.service
After=bluetooth.service