# usr/lib/systemd/system/user@.service.d/moeOS-hardening.conf
[Service]
KeyringMode=private
ProtectKernelModules=yes
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
#PrivateTmp=yes
CapabilityBoundingSet=~CAP_BPF
CapabilityBoundingSet=~CAP_SYS_TTY_CONFIG
CapabilityBoundingSet=~CAP_SYS_PTRACE
#ProtectProc=invisible
CapabilityBoundingSet=~CAP_NET_BIND_SERVICE
CapabilityBoundingSet=~CAP_SYS_ADMIN
