# usr/lib/systemd/system/user@.service.d/moeOS.conf
[Unit]
Wants=bluetooth.service
After=bluetooth.service