# usr/lib/systemd/system/user@.service.d/moeOS-hardening.conf
[Service]
KeyringMode=shared
#NoNewPrivileges=yes
RemoveIPC=yes
#ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProcSubset=pid
PrivateTmp=yes
LockPersonality=yes
