# usr/lib/systemd/system/cups.service.d/moeOS-resource-control.conf
[Service]
Nice=10
CPUWeight=20
IOWeight=50