# usr/lib/systemd/system/cups.service.d/moeOS-fix-stuck-printing.conf
[Service]
ExecStartPre=/usr/bin/find /var/spool/cups -maxdepth 1 -type f -delete