# usr/lib/systemd/system/user@service.service.d/override.conf
[Service]
ManagedOOMSwap=kill
ManagedOOMMemoryPressure=kill
