# usr/lib/systemd/system/devtools.slice
[Unit]
Description=moeOS Devtools Slice
Before=slices.target
