# usr/lib/systemd/system/devtools.slice
[Unit]
Description=moeOS Devtools Slice
Before=slices.target

[Slice]
MemoryHigh=80%
CPUWeight=idle
IOWeight=1
