# usr/lib/systemd/system/moeOS-noproxy.slice
[Unit]
Description=moeOS Proxy Exclude Slice
Before=slices.target
