# usr/lib/systemd/system/moeOS-scheduler.service
[Unit]
Description=moeOS scx_sched schedulers
StartLimitIntervalSec=30
StartLimitBurst=2
Conflicts=scx.service scx_loader.service

[Service]
Type=simple
#ExecStart=/usr/bin/scx_lavd --autopower
ExecStart=/usr/bin/scx_flash -m powersave -I 1000 -t 1000 -s 9000 -S 768
CPUWeight=10000
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target