# usr/lib/systemd/system/moeOS-scheduler.service
[Unit]
Description=moeOS scx_sched schedulers
StartLimitIntervalSec=30
StartLimitBurst=2
Conflicts=scx.service scx_loader.service

[Service]
Type=simple
#ExecStart=/usr/bin/scx_lavd --autopower
ExecStart=/usr/bin/scx_cosmos --no-deferred-wakeup --cpu-busy-thresh 50
CPUWeight=500
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target