# usr/lib/systemd/system/proxy.slice
[Unit]
Description=moeOS Proxy Slice
Before=slices.target
