# usr/lib/systemd/user/pipewire.service.d/moeOS-weight.conf
[Service]
CPUWeight=500
CPUSchedulingPriority=99
IOSchedulingPriority=0