# usr/lib/systemd/system/fstrim.service.d/moeOS-resource-control.conf
[Unit]
ConditionACPower=true

[Service]
MemoryHigh=30%
CPUQuota=50%
CPUWeight=idle
IOWeight=1
IOReadBandwidthMax=/ 50M