# usr/lib/systemd/system/tuned.service.d/moeOS-hardening.conf
[Service]
ProtectHome=tmpfs
BindPaths=/usr/share/moeOS-Docs/tuned:/etc/tuned
BindPaths=/usr/share/moeOS-Docs/tuned/active_profile:/etc/tuned/active_profile
BindPaths=/usr/share/moeOS-Docs/tuned/post_loaded_profile:/etc/tuned/post_loaded_profile
BindPaths=/usr/share/moeOS-Docs/tuned/bootcmdline:/etc/tuned/bootcmdline
PrivateTmp=disconnected
PrivateMounts=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
ProtectClock=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_NETLINK AF_UNIX
SocketBindDeny=ipv4:tcp
SocketBindDeny=ipv4:udp
SocketBindDeny=ipv6:tcp
SocketBindDeny=ipv6:udp