# etc/xdg/autostart/moeOS.desktop
[Desktop Entry]
Name=moeOS Setup
Exec=/usr/bin/moeOS-setup-user
Terminal=false
StartupNotify=false
Icon=moeos