# usr/lib/systemd/system/nvidia-powerd.service.d/moeOS-rc.conf
[Unit]
StartLimitIntervalSec=10
StartLimitBurst=9

[Service]
RestartSec=1s