# usr/lib/systemd/system/nvidia-powerd.service.d/moeOS-rc.conf
[Unit]
StartLimitIntervalSec=2
StartLimitBurst=2