# usr/lib/systemd/system/man-db.service.d/moeOS-rc.conf
[Unit]
ConditionACPower=true

[Service]
CPUWeight=idle
CPUQuota=5%