# usr/lib/systemd/system/geoclue.service.d/moeOS-hardening.conf
[Service]
PrivateTmp=disconnected
PrivateMounts=true
ProtectKernelLogs=true
ProtectHome=tmpfs