# usr/lib/systemd/system/power-profiles-daemon.service.d/moeOS-hardening.conf
[Service]
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=true
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectProc=ptraceable
MemoryDenyWriteExecute=yes
#RestrictAddressFamilies=AF_NETLINK AF_UNIX
LockPersonality=yes
RestrictRealtime=yes
ProtectClock=yes
SystemCallFilter=~@aio:EPERM @chown:EPERM @clock:EPERM @cpu-emulation:EPERM @debug:EPERM @ipc:EPERM @keyring:EPERM @memlock:EPERM @module:EPERM @mount:EPERM @obsolete:EPERM @pkey:EPERM @privileged:EPERM @raw-io:EPERM @reboot:EPERM @resources:EPERM @sandbox:EPERM @setuid:EPERM @swap:EPERM @timer:EPERM
