# usr/lib/systemd/user/localsearch-3.service.d/moeOS-resource-control.conf
[Service]
CPUQuota=75%
CPUWeight=idle
IOWeight=1
MemoryHigh=1000M
MemoryMax=1200M
Restart=always

[Unit]
OnFailure=moeOS-tracker-reset.service