# usr/lib/systemd/user/localsearch-3.service.d/moeOS-resource-control.conf
[Service]
CPUQuota=50%
CPUWeight=idle
IOWeight=1
MemoryHigh=100M
#MemoryMax=200M