# usr/lib/systemd/user/org.gnome.SettingsDaemon.Datetime.service.d/moeOS-override.conf
[Service]
Slice=background.slice
CPUQuota=1%
CPUWeight=idle
IOWeight=1
MemoryMax=120M
Restart=always
Nice=19

[Unit]
OnFailure=moeOS-tracker-reset.service